# etc/profile.d/chruby.sh
. /usr/local/etc/profile.d/chruby.sh

RUBIES=(
  /usr/local/ruby-1.8.7
  /usr/local/jruby-1.6.7.2
  /usr/local/jruby-1.7.0-preview2
  /usr/local/rubinius
)
