# etc/profile.d/chruby.sh
#!/bin/sh

source /usr/local/share/chruby/chruby.sh

RUBIES=(/opt/rubies/*)
